# backend/cron/crontab
# Run the backup script daily at 2 AM
0 2 * * * /usr/local/bin/db-backup.sh >> /var/log/cron/cron.log 2>&1
